# L3/MLMOS/1/csfncsf.service
[Unit]
Description=Serviciu Jmeker

[Service]
Type=simple
ExecStart=/bin/bash /root/boostrap.sh

[Install]
WantedBy=multi-user.target
